# service/xenproxy.service
[Unit]
Description=XenProxy LXC Multi-Protocol Gateway Admin Panel (SQLite)
After=network.target lxc.service
Requires=lxc.service

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__PYTHON_PATH__ __INSTALL_DIR__/app.py
Restart=on-failure
RestartSec=5
User=root
Group=root

# Ensure the application has write access to the database
# and can create network bridges
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_ADMIN
Environment="FLASK_ENV=production"
Environment="PRODUCTION=true"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xenproxy

[Install]
WantedBy=multi-user.target
