# service/xenproxy.service
[Unit]
Description=XenProxy LXC Multi-Protocol Gateway Admin Panel (SQLite)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/xenproxy
ExecStart=/usr/bin/python3 /root/xenproxy/app.py
Restart=on-failure
RestartSec=5
User=root
Group=root

# Ensure the application has write access to the database
# and can create network bridges
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_ADMIN
Environment="FLASK_ENV=production"

[Install]
WantedBy=multi-user.target
